# recipes/tedge-firmware-update/files/firmware-auto-rollback.timer
[Unit]
Description=rugpi auto rollback

[Timer]
OnBootSec=1800s

[Install]
WantedBy=timers.target